# Tools/Install/Machines/010_Xubuntu/Scripts/010_Basic_Config_/home/methods/.config/autostart/start_terminal.desktop
[Desktop Entry]
Name=xfce4-terminal
GenericName=Terminal
Comment=Interactive terminal
Exec=dbus-launch xfce4-terminal
Terminal=false
Type=Application
