# Tools/Install/Machines/010_Xubuntu/Scripts/010_Basic_Make-Root-Links_/etc/cron.d/run-at-first-boot
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
@reboot root touch /tmp/test-for-llorracc
@reboot root chmod a+rwx /tmp/test-for-llorracc
